# test_data/com.google.code.grr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.google.code.grr</string>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/lib/grr/grr_3.0.0.5_amd64/grr</string>
		<string>--config=/usr/lib/grr/grr_3.0.0.5_amd64/grr.yaml</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>ServiceDescription</key>
	<string>grr</string>
	<key>StandardOutPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
