# config/debian/dpkg_client/debian/grr-client.postinst.in
#!/bin/sh

set -e

DAEMON="%(ClientBuilder.target_dir)/%(Client.binary_name)"
DAEMON_ARGS="--config=%(ClientBuilder.target_dir)/%(ClientBuilder.config_filename)"

# This package is designed to work on systems with and without upstart.
case "$1" in
  configure\)
    ${DAEMON} ${DAEMON_ARGS} "--install"

    if [ -x /sbin/initctl ] && /sbin/initctl version | /bin/grep -q upstart; then
      # We have upstart, remove the init.d and any old rc links to prevent it
      # from being used.
      rm -f /etc/init.d/%(ClientBuilder.package_name)
      update-rc.d %(ClientBuilder.package_name) remove
      # Early versions of upstart didn't support restarting a service that
      # wasn't already running:
      # https://bugs.launchpad.net/ubuntu/+source/upstart/+bug/430883
      /usr/sbin/service %(ClientBuilder.package_name) stop 2>/dev/null || true
      /usr/sbin/service %(ClientBuilder.package_name) start 2>/dev/null
    else
      # Remove upstart config and use init.d
      rm -f /etc/init/%(ClientBuilder.package_name).conf
      /etc/init.d/%(ClientBuilder.package_name) restart 2>/dev/null
      # Update debian runlevel links, upstart handles this for ubuntu
      update-rc.d %(ClientBuilder.package_name) defaults
    fi
  ;;

  abort-upgrade|abort-remove|abort-deconfigure\)
  ;;

  *\)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
