# grr/core/install_data/macosx/client/grr.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>%(Client.plist_label)</string>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>%(ClientBuilder.install_dir)/%(Client.binary_name)</string>
		<string>--config=%(ClientBuilder.install_dir)/%(ClientBuilder.config_filename)</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>120</integer>
	<key>ServiceDescription</key>
	<string>%(Client.name)</string>
	<key>StandardOutPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
