# grr/core/install_data/debian/dpkg_client/upstart/grr-client.conf
# %(Client.description) client upstart file

limit nofile 20000 20000

kill timeout 300

start on startup
start on runlevel [2345]
stop on runlevel [016]

respawn

env LANG=en_US.UTF-8

script
	DAEMON="%(ClientBuilder.daemon_link).nanny"
	DAEMON_ARGS="%(ClientBuilder.daemon_link) --config=%(ClientBuilder.target_dir)/%(ClientBuilder.config_filename)"

	[ -x "${DAEMON}" ] || exit 0

	exec ${DAEMON} ${DAEMON_ARGS}

end script
