# config/centos/upstart/grr-http-server.conf
# GRR Rapid Response HTTP server upstart file

limit nofile 20000 20000

kill timeout 300

pre-start script
	mkdir -p /var/log/grr
	mkdir -p /var/run/grr/tmp/httpserver
end script

start on runlevel [2345]
stop on runlevel [016]

env MPLCONFIGDIR='/var/run/grr/tmp/httpserver'
env PYTHON_EGG_CACHE='/var/run/grr/tmp/httpserver'

script
	DAEMON="/usr/bin/grr_server"
	DEFAULTFILE="/etc/default/grr-http-server";

	[ -x "${DAEMON}" ] || DAEMON="/usr/local/bin/grr_server"

	[ -x "${DAEMON}" ] || exit 0

	if ! [ -f "${DEFAULTFILE}" ];
	then
		START="yes"
	else
		. ${DEFAULTFILE}
	fi

	if [ "x${START}" = "xyes" ];
	then
		exec ${DAEMON} --start_http_server --config=/etc/grr/grr-server.yaml
	fi
end script
